# port/files/files/etc/init.d/air720
#!/bin/sh /etc/rc.common

START=95
STOP=10

wait_boot_reset() {
	count=0
	while [ "$count" -lt 150 ]; do
		[ -e /tmp/air720-boot-reset.done ] && return 0
		count=$((count + 1))
		sleep 1
	done
	return 1
}

start() {
	mode="$(uci -q get air720.main.mode || echo hybrid)"

	# S25 air720-boot-reset owns the modem reset on normal boots. The first-boot
	# uci-defaults path also creates this marker after performing the same reset.
	wait_boot_reset || logger -t air720 "Boot reset completion marker timed out" 2>/dev/null || true

	if [ "$(cat /tmp/air720-control.state 2>/dev/null || echo idle)" = "hardware-reset-required" ]; then
		logger -t air720 "Air720SL AT channel is unresponsive; hardware reset required, skipping cellular bring-up" 2>/dev/null || true
		return 0
	fi

	case "$mode" in
		hybrid)
			/usr/bin/air720-control start >/tmp/air720-control-start.log 2>&1 &
			;;
		ppp)
			/usr/bin/air720-ppp start >/tmp/air720-ppp.log 2>&1 &
			;;
		rndis6)
			/usr/bin/air720-rndis6 start >/tmp/air720-rndis6.log 2>&1 &
			;;
		disabled|off|0)
			return 0
			;;
		*)
			echo "Unknown Air720 mode: $mode" >&2
			return 1
			;;
	esac
}

stop() {
	/usr/bin/air720-rndis6 stop >/dev/null 2>&1 || true
	/usr/bin/air720-ppp stop >/dev/null 2>&1 || true
}

restart() {
	stop
	sleep 2
	start
}
