# port/files/files/etc/hotplug.d/tty/95-air720-ppp
#!/bin/sh

VID="${AIR720_VID:-1286}"
PID="${AIR720_PID:-4e3d}"
mode="$(uci -q get air720.main.mode || echo hybrid)"
case "$mode" in hybrid|ppp) ;; *) exit 0 ;; esac

devname="${DEVNAME:-${DEVICENAME:-${DEVPATH##*/}}}"
[ -n "$devname" ] || exit 0
case "$devname" in /dev/*) dev="$devname" ;; *) dev="/dev/$devname" ;; esac

case "${ACTION:-}" in
	remove)
		configured="$(uci -q get air720.main.ppp_port || true)"
		[ "$configured" = "$dev" ] && ifdown cellular >/dev/null 2>&1 || true
		exit 0
		;;
	add) ;;
	*) exit 0 ;;
esac

[ -c "$dev" ] || exit 0
tty="${dev##*/}"
p="$(readlink -f "/sys/class/tty/$tty/device" 2>/dev/null || true)"
ifnum=''
match=0

while [ -n "$p" ] && [ "$p" != / ]; do
	[ -n "$ifnum" ] || ifnum="$(cat "$p/bInterfaceNumber" 2>/dev/null || true)"
	if [ -r "$p/idVendor" ] && [ -r "$p/idProduct" ]; then
		vid="$(tr 'A-F' 'a-f' < "$p/idVendor")"
		pid="$(tr 'A-F' 'a-f' < "$p/idProduct")"
		[ "$vid" = "$VID" ] && [ "$pid" = "$PID" ] && match=1
		break
	fi
	p="${p%/*}"; [ -n "$p" ] || p=/
done

[ "$match" -eq 1 ] || exit 0
case "$ifnum" in 03|3) ;; *) exit 0 ;; esac

(
	sleep 1
	case "$mode" in
		hybrid) /usr/bin/air720-control recover >/tmp/air720-control-tty-hotplug.log 2>&1 || true ;;
		ppp) /usr/bin/air720-ppp start >/tmp/air720-ppp-tty-hotplug.log 2>&1 || true ;;
	esac
) &
exit 0
