# port/files/files/etc/hotplug.d/net/95-air720-rndis
#!/bin/sh

# Ignore USB enumeration events until the persistent-power modem boot reset completes.
[ -e /tmp/air720-boot-reset.done ] || exit 0

[ -n "${INTERFACE:-}" ] || exit 0

case "${ACTION:-}" in
	remove)
		configured="$(uci -q get air720.main.rndis_device || true)"
		[ "$configured" = "$INTERFACE" ] && ifdown cellular6 >/dev/null 2>&1 || true
		exit 0
		;;
	add) ;;
	*) exit 0 ;;
esac

[ "$(uci -q get air720.main.rndis_ipv6 || echo 1)" = 1 ] || exit 0
[ -e "/sys/class/net/$INTERFACE/device/driver" ] || exit 0

driver="$(basename "$(readlink -f "/sys/class/net/$INTERFACE/device/driver")")"
case "$driver" in rndis_host|cdc_ether) ;; *) exit 0 ;; esac

p="$(readlink -f "/sys/class/net/$INTERFACE/device")"
match=0
while [ -n "$p" ] && [ "$p" != / ]; do
	if [ -r "$p/idVendor" ] && [ -r "$p/idProduct" ]; then
		vid="$(tr 'A-F' 'a-f' < "$p/idVendor")"
		pid="$(tr 'A-F' 'a-f' < "$p/idProduct")"
		[ "$vid" = 1286 ] && [ "$pid" = 4e3d ] && match=1
		break
	fi
	p="${p%/*}"; [ -n "$p" ] || p=/
done
[ "$match" -eq 1 ] || exit 0

mode="$(uci -q get air720.main.mode || echo hybrid)"
(
	sleep 1
	case "$mode" in
		hybrid) /usr/bin/air720-control recover >/tmp/air720-control-net-hotplug.log 2>&1 || true ;;
		rndis6) /usr/bin/air720-rndis6 start "$INTERFACE" >/tmp/air720-rndis6-hotplug.log 2>&1 || true ;;
	esac
) &
exit 0
