# port/files/files/etc/hotplug.d/net/95-air720-rndis
#!/bin/sh

[ "$ACTION" = "add" ] || exit 0
[ -n "$INTERFACE" ] || exit 0
[ -e "/sys/class/net/$INTERFACE/device/driver" ] || exit 0

driver="$(basename "$(readlink -f "/sys/class/net/$INTERFACE/device/driver")")"

case "$driver" in
	rndis_host|cdc_ether)
		;;
	*)
		exit 0
		;;
esac

# Dedicated Cellu-AR750S cellular interface. Wired WAN keeps the lower
# default metric, so it remains preferred when both uplinks are present.
uci -q set network.cellular='interface'
uci -q set network.cellular.device="$INTERFACE"
uci -q set network.cellular.proto='dhcp'
uci -q set network.cellular.metric='20'
uci -q commit network

# Put the cellular interface into the existing WAN firewall zone.
wan_zone="$(uci -q show firewall | sed -n "s/^\(firewall\.[^.]*\)\.name='wan'$/\1/p" | head -n1)"
if [ -n "$wan_zone" ]; then
	if ! uci -q get "${wan_zone}.network" 2>/dev/null | tr ' ' '\n' | grep -qx 'cellular'; then
		uci -q add_list "${wan_zone}.network=cellular"
		uci -q commit firewall
	fi
fi

ifup cellular >/dev/null 2>&1 &
exit 0
