# port/files/files/etc/uci-defaults/99-air720
#!/bin/sh

/etc/init.d/air720-boot-reset enable
/etc/init.d/air720 enable

# The S25 symlink is created only now on the first boot, after rcS enumeration.
# Perform the reset synchronously once before starting the data plane.
if [ "$(uci -q get air720.main.reset_on_boot || echo 1)" = 1 ]; then
	rm -f /tmp/air720-boot-reset.done
	/usr/bin/air720-control boot-reset >/tmp/air720-boot-reset-firstboot.log 2>&1 || true
fi

touch /tmp/air720-boot-reset.done
/etc/init.d/air720 start >/dev/null 2>&1 &
exit 0
