# port/files/files/etc/uci-defaults/99-air720
#!/bin/sh

/etc/init.d/air720-boot-reset enable
/etc/init.d/air720 enable

# On the first boot the S25 symlink is created while uci-defaults is already
# running. Execute the exact same reset service synchronously once. If rcS later
# discovers S25, its done marker makes it a no-op instead of resetting during PPP.
rm -f /tmp/air720-boot-reset.done
/etc/init.d/air720-boot-reset start >/tmp/air720-boot-reset-firstboot.log 2>&1 || true

/etc/init.d/air720 start >/dev/null 2>&1 &
exit 0
