# port/files/files/etc/uci-defaults/99-air720
#!/bin/sh

/etc/init.d/air720 enable
/etc/init.d/air720 start >/dev/null 2>&1 &
exit 0
